# tools/planepuck-enroll.service
[Unit]
Description=PlanePuck MQTT enrollment endpoint
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
# Runs as root so it can write the Mosquitto password file and reload the broker. The systemd
# sandbox below confines it (read-only FS except /etc/mosquitto, no home, private /tmp, no new privs).
# Locked-down alternative: run as a dedicated user, set PLANEPUCK_RELOAD_CMD="sudo systemctl reload
# mosquitto", grant a tight sudoers rule, and make the passwd file group-writable — see OTA-SETUP.md.
User=root
Environment=PLANEPUCK_PASSWD_FILE=/etc/mosquitto/passwd
# QUOTE the whole NAME=VALUE: systemd splits Environment= on spaces, so an UNQUOTED value here keeps
# only "systemctl" and drops "reload mosquitto" -> enroll-server.py runs bare `systemctl` (exits 0, a
# silent no-op), reports enroll success, but the broker never reloads the new password -> every login
# fails with MQTT state=5. The quotes preserve the full command.
Environment="PLANEPUCK_RELOAD_CMD=systemctl reload mosquitto"
Environment=PLANEPUCK_BIND=127.0.0.1:8090
# TOFU key-pin store (root-only; the service writes it atomically). Maps friend code -> per-device key,
# so after a code's first enroll only that device can change its password. Default matches the script.
Environment=PLANEPUCK_PIN_FILE=/etc/planepuck/enroll-keys.json
# Once the whole fleet runs signed-enroll firmware, uncomment to REJECT all unsigned/legacy enrolls:
# Environment=PLANEPUCK_REQUIRE_SIG=1
# The enroll token is the ONLY secret — keep it out of this unit; put it in a root-only env file:
#   /etc/planepuck/enroll.env   ->   PLANEPUCK_ENROLL_TOKEN=<same value as ENROLL_TOKEN in config.h>
#   chmod 600 /etc/planepuck/enroll.env
EnvironmentFile=/etc/planepuck/enroll.env
ExecStart=/usr/bin/python3 /opt/planepuck/enroll-server.py
Restart=on-failure
RestartSec=2

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/etc/mosquitto /etc/planepuck   # passwd file + the TOFU pin store (ProtectSystem=full else RO)
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
